# packaging/debian/DEBIAN/prerm
#!/bin/sh
set -e

# stop/disable on remove
if [ "$1" = "remove" ] || [ "$1" = "deconfigure" ]; then
    if command -v systemctl >/dev/null 2>&1; then
        systemctl stop ngm.service || true
        systemctl disable ngm.service || true
    fi
fi

exit 0
