# configs/ngm.service
[Unit]
Description=ngm
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/ngm
ExecStart=/opt/ngm/bin/ngm
Restart=always
RestartSec=3
KillSignal=SIGTERM
TimeoutStopSec=30

# Optional env file
EnvironmentFile=-/etc/default/ngm

# Basic hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true

[Install]
WantedBy=multi-user.target
